# Utilities/skimage_watchdog.service
[Unit]
Description=Starts Skimage watchdog
After=docker.service

[Service]
ExecStart=xfce4-terminal

[Install]
WantedBy=multi-user.target
